# config/messenger-worker.ini
[program:messenger-consume]
command=php81 -d memory_limit=12G /usr/home/pvpluy/public_html/dev.pvp/bin/console messenger:consume async --limit=10 --time-limit=14400
user=pvpluy
numprocs=2
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
